# files/systemd/mysite.socket
[Unit]
Description=mysite socket

[Socket]
ListenStream=/run/mysite/mysite.sock

[Install]
WantedBy=sockets.target